# sway/dot-config/systemd/user/kmonad.service
[Unit]
Description=KMonad keyboard config

[Service]
Restart=always
RestartSec=3
ExecStart=%h/.local/bin/kmonad %h/.config/config.kbd
Nice=-20
